# bulletin.service
[Unit]
Description=교회 모바일 주보 시스템
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/bulletin
EnvironmentFile=/opt/bulletin/.env
Environment="SPRING_PROFILES_ACTIVE=ec2"
ExecStart=/usr/bin/java -Xms512m -Xmx1024m -XX:+UseG1GC -XX:MaxGCPauseMillis=200 -jar /opt/bulletin/bulletin-0.0.1-SNAPSHOT.jar
ExecStop=/bin/kill -15 $MAINPID
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bulletin

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bulletin/uploads /opt/bulletin/logs

# 리소스 제한
LimitNOFILE=65536
LimitNPROC=4096
MemoryLimit=2G

[Install]
WantedBy=multi-user.target

